# files/etc/uci-defaults/zzz-010-lff-system
#!/bin/sh

local macaddr="$(uci get wireless.radio0.macaddr)"

uci batch <<EOF
set system.@system[0].hostname='Freifunk-${macaddr//:/}'

commit system
EOF
